# launchagents/com.welisa.callbridge-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Replaces the legacy com.autologger.server.plist as part of BUILD-04 (unified com.welisa.* namespace). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.welisa.callbridge-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/Caskroom/miniconda/base/bin/python3.13</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>main:app</string>
        <string>--host</string>
        <string>localhost</string>
        <string>--port</string>
        <string>8765</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kiranknoppert/dev/callbridge</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>logs/uvicorn.out.log</string>
    <key>StandardErrorPath</key>
    <string>logs/uvicorn.err.log</string>
</dict>
</plist>
